# templates/etc/init/ttrss-update-daemon.conf
description	"ttrss-update-daemon"

start on runlevel [2345]
stop on runlevel [!2345]

setuid ttrss

exec /usr/bin/php /srv/ttrss/update_daemon2.php

